# packages/ambra-cas/src/main/scripts/ambra-cas
#!/usr/bin/env bash
# $HeadURL::                                                                   $
# $Id$
#
# Tomcat init script
# Copyright (c) 2006-2008 by Topaz, Inc.
# http://topazproject.org
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

# chkconfig: 2345 91 15
# description: CAS authentication service and plos-registration

# Source function library (not really used but seems standard on redhat)
. /etc/init.d/functions

# Read our configuration from sysconfig (a redhat convention)
[ -e /etc/sysconfig/topaz-common ] && . /etc/sysconfig/topaz-common
[ -e /etc/sysconfig/ambra-cas ] && . /etc/sysconfig/ambra-cas

# If our system is setup incorrectly somehow, minimal attempt to clean up
[ -n "$JAVA_HOME" ] && [ ! -d "$JAVA_HOME" ] && \
    echo "WARN: Resetting Invalid JAVA_HOME ($JAVA_HOME)" && JAVA_HOME=
    
# Try to set JAVA_HOME if not set above
[ -z "$JAVA_HOME" ] && JAVA_HOME=`find /opt -maxdepth 1 -name 'jdk*' -a \! -type f | head -1`
export JAVA_HOME

# Start/stop/etc tomcat
/usr/local/ambra-cas/bin/tomcat --basename $(basename `readlink -f $0`) "$@"
