# packages/plosone/src/main/scripts/preremove
if [ $1 = 0 ]; then
    service plosone stop > /dev/null 2>&1
    /sbin/chkconfig --del plosone || :
fi
