# packages/mulgara/src/main/scripts/mulgara
#!/usr/bin/env bash
# $HeadURL::                                                                   $
# $Id$
#
# Tomcat init script
# Copyright (c) 2006 by Topaz, Inc.
# http://topazproject.org
#
# Licensed under the Educational Community License version 1.0
# http://opensource.org/licenses/ecl1.php

# chkconfig: 2345 92 14
# description: mulgara rdf database

# Validate valid configuration files
[ -e /etc/topaz/topaz.xml ] && ! xmllint --noout /etc/topaz/topaz.xml && exit 1

# Source function library (not really used but seems standard on redhat)
. /etc/init.d/functions

# Read our configuration from sysconfig (a redhat convention)
[ -e /etc/sysconfig/topaz-common ] && . /etc/sysconfig/topaz-common
[ -e /etc/sysconfig/mulgara ] && . /etc/sysconfig/mulgara

# If our system is setup incorrectly somehow, minimal attempt to clean up
[ -n "$JAVA_HOME" ] && [ ! -d "$JAVA_HOME" ] && \
    echo "WARN: Resetting Invalid JAVA_HOME ($JAVA_HOME)" && JAVA_HOME=
    
# Try to set JAVA_HOME if not set above
[ -z "$JAVA_HOME" ] && JAVA_HOME=`find /opt -name 'jdk*' -maxdepth 1 | head -1`
export JAVA_HOME

# Start multicast for ecache
/usr/local/topaz/bin/multicast

# Start/stop/etc tomcat
/usr/local/topaz/bin/tomcat --basename $(basename `readlink -f $0`) "$@"
